# supervisord-local.conf
[supervisord]
    logfile = /tmp/supervisord.log
    loglevel = info
    pidfile = /tmp/supervisord.pid
    nodaemon = true
    directory = /tmp
    nocleanup = true
    childlogdir = /tmp
    user=root

; BR24 Download Upload Server supervisor config file.
; want to test out without the command to kill all the processes if one of them fails to start

[eventlistener:processes]
    command=bash -c "printf 'READY\n' && while read line; do kill -SIGQUIT $PPID; done < /dev/stdin"
    events=PROCESS_STATE_STOPPED,PROCESS_STATE_EXITED,PROCESS_STATE_FATAL

[program:php8-fpm]
    priority=100
    command=php-fpm -F --nodaemonize --fpm-config=/usr/local/etc/php-fpm.d/www.conf
    stdout_logfile=/dev/stdout
    stdout_logfile_maxbytes=0
    stderr_logfile=/dev/stderr
    stderr_logfile_maxbytes=0
    exitcodes=0

[program:crond-service]
    priority=300
    command=/usr/sbin/cron -f
    stdout_logfile=/var/www/src/alpine/storage/logs/openrc/crond.log
    stderr_logfile=/var/www/src/alpine/storage/logs/openrc/crond.err
    user=root


; BR24 Download Upload Server supervisor config file [LARAVEL WORKERS].
; the websockets supervisord command always fails when run on docker local
;[program:websockets-service]
;    priority=310
;    command=/usr/local/bin/php /var/www/src/alpine/artisan websockets:serve --port=6008
;    stdout_logfile=/var/www/src/alpine/storage/logs/openrc/websockets.log
;    stderr_logfile=/var/www/src/alpine/storage/logs/openrc/websockets.err
;    user=www-data
;    numprocs=1
;    startretries=3
;    startsecs=10
;    autostart=true
;    autorestart=unexpected
;    exitcodes=0

[program:queue-worker-service]
    priority=400
    command=/usr/local/bin/php /var/www/src/alpine/artisan queue:work --queue=default
    stdout_logfile=/var/www/src/alpine/storage/logs/openrc/default.log
    stderr_logfile=/var/www/src/alpine/storage/logs/openrc/default.err
    user=www-data
    numprocs=1
    startretries=3
    startsecs=10
    autostart=true
    autorestart=unexpected
    exitcodes=0

[program:queue-autodl_download]
    priority=410
    command=/usr/local/bin/php /var/www/src/alpine/artisan queue:work --queue=autodl_download
    stdout_logfile=/var/www/src/alpine/storage/logs/openrc/autodl_download.log
    stderr_logfile=/var/www/src/alpine/storage/logs/openrc/autodl_download.err
    user=www-data
    numprocs=1
    startretries=3
    startsecs=10
    autostart=true
    autorestart=unexpected
    exitcodes=0

[program:queue-autodl_extractedcheckscan]
    priority=411
    command=/usr/local/bin/php /var/www/src/alpine/artisan queue:work --queue=autodl_extractedcheckscan
    stdout_logfile=/var/www/src/alpine/storage/logs/openrc/autodl_extractedcheckscan.log
    stderr_logfile=/var/www/src/alpine/storage/logs/openrc/autodl_extractedcheckscan.err
    user=www-data
    numprocs=1
    startretries=3
    startsecs=10
    autostart=true
    autorestart=unexpected
    exitcodes=0

[program:queue-autodl_messageaftermovetodirectory]
    priority=412
    command=/usr/local/bin/php /var/www/src/alpine/artisan queue:work --queue=autodl_messageaftermovetodirectory
    stdout_logfile=/var/www/src/alpine/storage/logs/openrc/autodl_messageaftermovetodirectory.log
    stderr_logfile=/var/www/src/alpine/storage/logs/openrc/autodl_messageaftermovetodirectory.err
    user=www-data
    numprocs=1
    startretries=3
    startsecs=10
    autostart=true
    autorestart=unexpected
    exitcodes=0

[program:queue-autodl_movetodirectory]
    priority=413
    command=/usr/local/bin/php /var/www/src/alpine/artisan queue:work --queue=autodl_movetodirectory
    stdout_logfile=/var/www/src/alpine/storage/logs/openrc/autodl_movetodirectory.log
    stderr_logfile=/var/www/src/alpine/storage/logs/openrc/autodl_movetodirectory.err
    user=www-data
    numprocs=1
    startretries=3
    startsecs=10
    autostart=true
    autorestart=unexpected
    exitcodes=0

[program:queue-autodl_unzip]
    priority=414
    command=/usr/local/bin/php /var/www/src/alpine/artisan queue:work --queue=autodl_unzip
    stdout_logfile=/var/www/src/alpine/storage/logs/openrc/autodl_unzip.log
    stderr_logfile=/var/www/src/alpine/storage/logs/openrc/autodl_unzip.err
    user=www-data
    numprocs=1
    startretries=3
    startsecs=10
    autostart=true
    autorestart=unexpected
    exitcodes=0

[program:queue-autodl_unzipchecks]
    priority=415
    command=/usr/local/bin/php /var/www/src/alpine/artisan queue:work --queue=autodl_unzipchecks
    stdout_logfile=/var/www/src/alpine/storage/logs/openrc/autodl_unzipchecks.log
    stderr_logfile=/var/www/src/alpine/storage/logs/openrc/autodl_unzipchecks.err
    user=www-data
    numprocs=1
    startretries=3
    startsecs=10
    autostart=true
    autorestart=unexpected
    exitcodes=0

[program:queue-manualdl_download]
    priority=420
    command=/usr/local/bin/php /var/www/src/alpine/artisan queue:work --queue=manualdl_download
    stdout_logfile=/var/www/src/alpine/storage/logs/openrc/manualdl_download.log
    stderr_logfile=/var/www/src/alpine/storage/logs/openrc/manualdl_download.err
    user=www-data
    numprocs=1
    startretries=3
    startsecs=10
    autostart=true
    autorestart=unexpected
    exitcodes=0

[program:queue-manualdl_extractedcheckscan]
    priority=421
    command=/usr/local/bin/php /var/www/src/alpine/artisan queue:work --queue=manualdl_extractedcheckscan
    stdout_logfile=/var/www/src/alpine/storage/logs/openrc/manualdl_extractedcheckscan.log
    stderr_logfile=/var/www/src/alpine/storage/logs/openrc/manualdl_extractedcheckscan.err
    user=www-data
    numprocs=1
    startretries=3
    startsecs=10
    autostart=true
    autorestart=unexpected
    exitcodes=0

[program:queue-manualdl_messageaftermovetodirectory]
    priority=422
    command=/usr/local/bin/php /var/www/src/alpine/artisan queue:work --queue=manualdl_messageaftermovetodirectory
    stdout_logfile=/var/www/src/alpine/storage/logs/openrc/manualdl_messageaftermovetodirectory.log
    stderr_logfile=/var/www/src/alpine/storage/logs/openrc/manualdl_messageaftermovetodirectory.err
    user=www-data
    numprocs=1
    startretries=3
    startsecs=10
    autostart=true
    autorestart=unexpected
    exitcodes=0

[program:queue-manualdl_movetodirectory]
    priority=423
    command=/usr/local/bin/php /var/www/src/alpine/artisan queue:work --queue=manualdl_movetodirectory
    stdout_logfile=/var/www/src/alpine/storage/logs/openrc/manualdl_movetodirectory.log
    stderr_logfile=/var/www/src/alpine/storage/logs/openrc/manualdl_movetodirectory.err
    user=www-data
    numprocs=1
    startretries=3
    startsecs=10
    autostart=true
    autorestart=unexpected
    exitcodes=0

[program:queue-manualdl_unzip]
    priority=424
    command=/usr/local/bin/php /var/www/src/alpine/artisan queue:work --queue=manualdl_unzip
    stdout_logfile=/var/www/src/alpine/storage/logs/openrc/manualdl_unzip.log
    stderr_logfile=/var/www/src/alpine/storage/logs/openrc/manualdl_unzip.err
    user=www-data
    numprocs=1
    startretries=3
    startsecs=10
    autostart=true
    autorestart=unexpected
    exitcodes=0

[program:queue-manualdl_unzipchecks]
    priority=425
    command=/usr/local/bin/php /var/www/src/alpine/artisan queue:work --queue=manualdl_unzipchecks
    stdout_logfile=/var/www/src/alpine/storage/logs/openrc/manualdl_unzipchecks.log
    stderr_logfile=/var/www/src/alpine/storage/logs/openrc/manualdl_unzipchecks.err
    user=www-data
    numprocs=1
    startretries=3
    startsecs=10
    autostart=true
    autorestart=unexpected
    exitcodes=0

[program:queue-manualul_afterzip_movetojobfolderreadydirectory]
    priority=430
    command=/usr/local/bin/php /var/www/src/alpine/artisan queue:work --queue=manualul_afterzip_movetojobfolderreadydirectory
    stdout_logfile=/var/www/src/alpine/storage/logs/openrc/manualul_afterzip_movetojobfolderreadydirectory.log
    stderr_logfile=/var/www/src/alpine/storage/logs/openrc/manualul_afterzip_movetojobfolderreadydirectory.err
    user=www-data
    numprocs=1
    startretries=3
    startsecs=10
    autostart=true
    autorestart=unexpected
    exitcodes=0

[program:queue-manualul_afterzip_movetojobfolderreadydirectory_send_message]
    priority=431
    command=/usr/local/bin/php /var/www/src/alpine/artisan queue:work --queue=manualul_afterzip_movetojobfolderreadydirectory_send_message
    stdout_logfile=/var/www/src/alpine/storage/logs/openrc/manualul_afterzip_movetojobfolderreadydirectory_send_message.log
    stderr_logfile=/var/www/src/alpine/storage/logs/openrc/manualul_afterzip_movetojobfolderreadydirectory_send_message.err
    user=www-data
    numprocs=1
    startretries=3
    startsecs=10
    autostart=true
    autorestart=unexpected
    exitcodes=0

[program:queue-manualul_check_the_progress_of_ready_zip_to_s3]
    priority=432
    command=/usr/local/bin/php /var/www/src/alpine/artisan queue:work --queue=manualul_check_the_progress_of_ready_zip_to_s3
    stdout_logfile=/var/www/src/alpine/storage/logs/openrc/manualul_check_the_progress_of_ready_zip_to_s3.log
    stderr_logfile=/var/www/src/alpine/storage/logs/openrc/manualul_check_the_progress_of_ready_zip_to_s3.err
    user=www-data
    numprocs=1
    startretries=3
    startsecs=10
    autostart=true
    autorestart=unexpected
    exitcodes=0

[program:queue-manualul_checkthenzip]
    priority=433
    command=/usr/local/bin/php /var/www/src/alpine/artisan queue:work --queue=manualul_checkthenzip
    stdout_logfile=/var/www/src/alpine/storage/logs/openrc/manualul_checkthenzip.log
    stderr_logfile=/var/www/src/alpine/storage/logs/openrc/manualul_checkthenzip.err
    user=www-data
    numprocs=1
    startretries=3
    startsecs=10
    autostart=true
    autorestart=unexpected
    exitcodes=0

[program:queue-manualul_send_the_ready_zip_to_s3]
    priority=434
    command=/usr/local/bin/php /var/www/src/alpine/artisan queue:work --queue=manualul_send_the_ready_zip_to_s3
    stdout_logfile=/var/www/src/alpine/storage/logs/openrc/manualul_send_the_ready_zip_to_s3.log
    stderr_logfile=/var/www/src/alpine/storage/logs/openrc/manualul_send_the_ready_zip_to_s3.err
    user=www-data
    numprocs=1
    startretries=3
    startsecs=10
    autostart=true
    autorestart=unexpected
    exitcodes=0
