# openrc_queue_scripts/automatic_download/laravel-queue-worker-service-extractedcheckscan
#!/sbin/openrc-run

# Laravel Queue Worker serve using openrc
# ----------------------------------
# https://wiki.gentoo.org/wiki/OpenRC_to_systemd_Cheatsheet
# https://wiki.gentoo.org/wiki/OpenRC/supervise-daemon
# place this file /etc/init.d/laravel-queue-worker-service-extractedcheckscan
# https://confluence.atlassian.com/kb/starting-service-on-linux-throws-a-no-such-file-or-directory-error-794203722.html
# use command (to remove carriage retun characters):- sudo sed -i -e 's/\r//g' /etc/init.d/laravel-queue-worker-service-extractedcheckscan
# make this file executable using:- sudo chmod a+x /etc/init.d/laravel-queue-worker-service-extractedcheckscan
# make the shell script executable:- sudo chmod a+x /var/www/src/alpine/openrc_queue_scripts/automatic_download/laravel-queue-worker-service-extractedcheckscan.sh
# run this command to enable service:
# sudo rc-update add laravel-queue-worker-service-extractedcheckscan
# sudo rc-service laravel-queue-worker-service-extractedcheckscan start
# sudo rc-service laravel-queue-worker-service-extractedcheckscan status

name=${RC_SVCNAME}
description="Laravel Queue Worker Serve"
supervisor="supervise-daemon"
command="/bin/sh"
command_args="/var/www/src/alpine/openrc_queue_scripts/automatic_download/${RC_SVCNAME}.sh"
command_user="root:www-data"
pidfile="/var/run/${RC_SVCNAME}.pid"
output_log="/var/www/src/alpine/storage/logs/openrc/$(date +%Y_%m_%d)_${RC_SVCNAME}.log"
error_log="/var/www/src/alpine/storage/logs/openrc/$(date +%Y_%m_%d)_${RC_SVCNAME}.err"

depend() {
    need net
}

reload() {
        ebegin "Reloading laravel-queue-worker-service-extractedcheckscan configuration"
#       start-stop-daemon --exec $command --signal HUP
        ${supervisor} ${RC_SVCNAME} --signal HUP --pidfile "${pidfile}"
        eend $?
}